# src/openvpn/.service
[Unit]
Requires=networking.service
After=networking.service
